# start_swift.service
[Unit]
Description=Mount swift drives at startup

[Service]
ExecStart=/opt/swift/bin/mount_devices.sh

[Install]
WantedBy=multi-user.target
